# homebox/rootfs/etc/s6-overlay/s6-rc.d/homebox/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Homebox
# Long-running service: starts the Homebox binary with the correct environment.
# ==============================================================================

# Fix: Homebox tries to create temp files under $TMPDIR.
# /tmp is not writable in HA addon containers, so we redirect to /data/tmp.
mkdir -p /data/tmp
chmod 1777 /data/tmp
export TMPDIR="/data/tmp"

# Change to /data so Homebox resolves all relative storage paths
# inside the addon's persistent data directory.
cd /data

# Use Homebox's official SQLite configuration variable and keep DB under
# persistent /data. Keep WAL mode enabled for SQLite durability.
export HBOX_DATABASE_SQLITE_PATH="/data/homebox.db?_pragma=busy_timeout=30000&_pragma=journal_mode=WAL&_fk=1&_time_format=sqlite"

# Show the effective SQLite database file path in startup logs.
HBOX_DATABASE_SQLITE_FILE="${HBOX_DATABASE_SQLITE_PATH%%\?*}"
mkdir -p "$(dirname "${HBOX_DATABASE_SQLITE_FILE}")"
bashio::log.info "Using SQLite database at ${HBOX_DATABASE_SQLITE_FILE}"

# Run in production mode (disables developer tooling/debug endpoints).
export HBOX_MODE="production"

# Apply user-configurable options from the HA add-on UI.
export HBOX_LOG_LEVEL="$(bashio::config 'log_level')"
export HBOX_LOG_FORMAT="text"
export HBOX_WEB_MAX_UPLOAD_SIZE="$(bashio::config 'max_upload_size')"
export HBOX_OPTIONS_ALLOW_REGISTRATION="$(bashio::config 'allow_registration')"

# Disable telemetry and self-update checks (not useful in a local HA context).
export HBOX_OPTIONS_ALLOW_ANALYTICS="false"
export HBOX_OPTIONS_GITHUB_RELEASE_CHECK="false"

# Trust reverse-proxy headers if HA is fronted by a proxy, but keep false
# for the default direct-access scenario.
export HBOX_OPTIONS_TRUST_PROXY="false"

if bashio::config.true 'ha_api_startup_check'; then
       if [ -n "${SUPERVISOR_TOKEN:-}" ]; then
	       HA_API_STATUS="$(curl -sS -o /dev/null -w "%{http_code}" -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" "http://supervisor/core/api/" || true)"
	       if [ "${HA_API_STATUS}" = "200" ]; then
		       bashio::log.info "Home Assistant API startup check succeeded"
	       else
		       bashio::log.warning "Home Assistant API startup check returned HTTP ${HA_API_STATUS}"
	       fi
       else
	       bashio::log.warning "SUPERVISOR_TOKEN not available; skipping Home Assistant API startup check"
       fi
fi

bashio::log.info "Starting Homebox..."
exec /usr/bin/homebox
